# dot_config/systemd/user/litellm.service
[Unit]
Description=LiteLLM Proxy
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/ianpascoe/.litellm
ExecStart=/home/ianpascoe/.litellm/venv/bin/litellm --config config.yaml
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
Environment=HOME=/home/ianpascoe
Environment=TMPDIR=/tmp
Environment=PATH=/home/ianpascoe/.litellm/venv/bin:/home/ianpascoe/.local/bin:/home/ianpascoe/.cache/.bun/bin:/home/ianpascoe/.local/share/pnpm:/home/ianpascoe/.npm-global/bin:/home/ianpascoe/.cargo/bin:/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin

[Install]
WantedBy=default.target
